# examples/docker-gen.upstart.conf
#!upstart

description "docker-gen upstart script"

start on startup
stop on shutdown

respawn

exec start-stop-daemon --start --exec /usr/bin/docker-gen -- -config /etc/docker-gen.cfg